# .docker/glpi-app/files/etc/cron.d/glpi
* * * * *    /usr/local/bin/php /var/www/glpi/front/cron.php &>/dev/null
